# network-config/zigbee-lan/startup-task.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Tarea del Programador de Windows que arranca automáticamente al iniciar el PC:
    1. Adjunta el Sonoff Dongle Max al WSL2 mediante usbipd
    2. Arranca ser2net en WSL2 para exponer el dongle en TCP:20108

  El script setup-hp-mini.ps1 rellena BUSID_PLACEHOLDER, DISTRO_PLACEHOLDER
  y PORT_PLACEHOLDER antes de registrar esta tarea.
  Si necesitas editar la tarea a mano, abre el Programador de tareas → ZigbeeUSBLAN.
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Conecta el Sonoff Dongle Max ZigBee al WSL2 y arranca ser2net para que Home Assistant pueda acceder al coordinador ZigBee por red.</Description>
    <Author>setup-hp-mini.ps1</Author>
  </RegistrationInfo>
  <Triggers>
    <BootTrigger>
      <Enabled>true</Enabled>
      <!-- Espera 45 s para que WSL2 y los servicios de red estén listos -->
      <Delay>PT45S</Delay>
    </BootTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <!-- SYSTEM tiene permisos para usbipd y no requiere sesión de usuario abierta -->
      <UserId>S-1-5-18</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <RestartOnFailure>
      <Interval>PT1M</Interval>
      <Count>3</Count>
    </RestartOnFailure>
  </Settings>
  <Actions Context="Author">
    <!-- Paso 1: adjuntar el USB al WSL2 -->
    <Exec>
      <Command>usbipd</Command>
      <Arguments>attach --wsl --busid BUSID_PLACEHOLDER</Arguments>
    </Exec>
    <!-- Paso 2: arrancar ser2net dentro de WSL2 en segundo plano -->
    <Exec>
      <Command>wsl</Command>
      <Arguments>-d DISTRO_PLACEHOLDER -- bash -c "sudo ser2net -c /etc/ser2net.yaml -n &amp;"</Arguments>
    </Exec>
  </Actions>
</Task>
